# meta-hassio/recipes-support/resinhup/resinhup/resinhup.timer
[Unit]
Description=Resinhup updater timer

[Timer]
OnBootSec=15min
OnUnitInactiveSec=1d

[Install]
WantedBy=resin.target
