# rootfs/etc/services.d/watchdog/finish
#!/usr/bin/env bashio
# ==============================================================================
# Take down the S6 supervision tree when Watchdog fails
# ==============================================================================

if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "Halt Supervisor (Wuff)"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Watchdog restart after closing"
