# rootfs/etc/services.d/supervisor/finish
#!/usr/bin/env bashio
# ==============================================================================
# Take down the S6 supervision tree when Supervisor fails
# ==============================================================================

if [[ "$1" -ne 100 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "Halt Supervisor"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Supervisor restart after closing"
