# rootfs/etc/services.d/supervisor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Service service
# ==============================================================================

# Enable Jemalloc for Home Assistant Core, unless disabled
if [[ -z "${DISABLE_JEMALLOC+x}" ]]; then
  export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"
fi
exec python3 -m supervisor
