# rootfs/etc/services.d/supervisor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Supervisor service
# ==============================================================================
export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"

exec python3 -m supervisor
